# sys-power/cpufrequtils/files/cpufrequtils-init.d-006
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

change() {
	ebegin "Running cpufreq-set ${opts}"
	/usr/libexec/cpufrequtils-change.sh "${@}"
	eend $?
}

start() {
	change "${START_OPTS}" ${SYSFS_EXTRA}
}

stop() {
	change "${STOP_OPTS}"
}
